# package/busybox/files/genoa/etc/inittab
ttyS3::sysinit:/etc/rc.local
ttyS3::respawn:-/bin/sh
::restart:/sbin/init
::ctrlaltdel:/bin/umount -a -r

