# package/busybox/files/mst12/etc/inittab
::sysinit:/etc/rc.local
::respawn:/bin/cttyhack -/bin/sh
::sysinit:/etc/rc.local
::respawn:/bin/cttyhack -/bin/sh
